# sys-power/asusctl/files/openrc/asusd-user
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors

description="ASUS User Daemon"

command="/usr/bin/asusd-user"
command_background=true
start_stop_daemon_args="-u ${USER}"

pidfile="/run/asuslinux/${RC_SVCNAME}.pid"
output_log="/var/log/asuslinux/${RC_SVCNAME}.log"

depend() {
    need asusd
}

checkconfig() {
    if [ "${RC_VERSION:-0}" = "0" ]; then
        eerror "This script cannot be used for baselayout-1."
        return 1
    fi

    if [ "${USER}" = "${RC_SVCNAME}" ]; then
        eerror "You have to create an init script for each user:"
	    eerror "ln -s asusd-user /etc/init.d/asusd-user.<user>"
	    return 1
    fi

    if ! id -u "${USER}" >/dev/null; then
        eerror "${USER}: No such user"
        return 1
    fi

    checkpath -f --owner "${USER}" --mode 0755 "${output_log}"
}

start_pre() {
    checkconfig || return $?
    sleep 2
}
