# net-proxy/clash-premium-bin/files/clash
#!/sbin/openrc-run
# Copyright 2008-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2 or later

: ${CLASH_CWD:=/usr/share/clash}

USER=${RC_SVCNAME##*.}
PIDFILE_DIR=/var/run/clash/${USER}
PIDFILE=${PIDFILE_DIR}/clash.pid
CLASH_LOG_DIR="/var/log/clash/${USER}"
CLASH_LOG_FILE="${CLASH_LOG_DIR}/clash-linux.log"

U_HOME="/home/${USER}"
U_CONF_DIR="${U_HOME}/.config/clash"


description="Clash service daemon."

depend() {
    need localmount
    after network-online nftables iptabels local
}

checkconfig() {
    if [ "${RC_VERSION:-0}" = "0" ]; then
        eerror "This script cannot be used for baselayout-1."
        return 1
    fi

    if [ "${USER}" = "${RC_SVCNAME}" ]; then
        eerror "You have to create an init script for each user:"
	    eerror "ln -s clash /etc/init.d/clash.<user>"
	    return 1
    fi

    if ! id -u "${USER}" >/dev/null; then
        eerror "${USER}: No such user"
        return 1
    fi

    if [ ! -f "${U_CONF_DIR}/config.yaml" ]; then
        eerror "No such config file: ${U_CONF_DIR}/config.yaml"
        return 1
    fi

    checkpath -d --owner 0 --mode 0755 "${PIDFILE_DIR%/*}"
    checkpath -d --owner "${USER}" --mode 0755 "${PIDFILE_DIR}"

    checkpath -d --owner 0 --mode 0755 "${CLASH_LOG_DIR%/*}"
    checkpath -d --owner "${USER}" --mode 0755 "${CLASH_LOG_DIR}"
}

start_pre() {
    checkconfig || return $?
}

start() {
    ebegin "Starting clash daemon for user ${USER}"

    local clash_opts=""
    local enable_tun=$(yq .tun.enable ${U_CONF_DIR}/config.yaml)

    if [ "${enable_tun}" == "true" ] ; then
       clash_opts="/usr/bin/clash-linux -d ${U_CONF_DIR}"

       sh ${CLASH_CWD}/setup-cgroup.sh > /dev/null

       start-stop-daemon --start --background \
            --pidfile "${PIDFILE}" --make-pidfile  \
            --stdout "${CLASH_LOG_FILE}" --stderr "${CLASH_LOG_FILE}" \
            --chdir "${U_HOME}" \
            --exec ${CLASH_CWD}/bypass-proxy -- ${clash_opts}
    else
       clash_opts="-d ${U_CONF_DIR}"

       start-stop-daemon --start --background \
            --pidfile "${PIDFILE}" --make-pidfile  \
            --stdout "${CLASH_LOG_FILE}" --stderr "${CLASH_LOG_FILE}" \
            --chdir "${U_HOME}" \
            --exec /usr/bin/clash-linux -- ${clash_opts}
    fi

    eend $?
}

stop() {
    ebegin "Stopping clash daemon for user ${USER}"
        start-stop-daemon --stop --pidfile "${PIDFILE}"
    eend $?
}
