# sys-power/supergfxctl/files/openrc/supergfxd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors

description="Supergfx Daemon"

command="/usr/bin/supergfxd"
command_background=true

pidfile="/run/asuslinux/${RC_SVCNAME}.pid"

depend() {
    need dbus
    after local
}

start_pre() {
    checkpath -q -d "/run/asuslinux"
}
