# sys-power/power-profiles-daemon/files/openrc/power-profiles-daemon
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors

description="Power Profiles daemon"

command="/usr/libexec/power-profiles-daemon"
command_background=true
#Uncomment this to enable debug
#start_stop_daemon_args="-e G_MESSAGES_DEBUG=all"

pidfile="/run/asuslinux/${RC_SVCNAME}.pid"

depend() {
    need dbus
    after local
}

start_pre() {
    checkpath -q -d "/run/asuslinux"
}
