# sys-power/auto-cpufreq/files/auto-cpufreq
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors

description="auto-cpufreq - Automatic CPU speed & power optimizer for Linux"

supervisor="supervise-daemon"

command="auto-cpufreq"
command_background=true
command_args="--daemon"
command_user="root"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}/${RC_SVCNAME}.log"

depend() {
    after net local
}

start_pre() {
    checkpath -q -d "/var/log/${RC_SVCNAME}"
    checkpath -f --owner "${USER}" --mode 0755 "${output_log}"
}
