# sys-power/asusctl/files/asusd-user
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors

description="ASUS User Daemon"
command="/usr/bin/asusd-user"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
supervisor=supervise-daemon

depend() {
        after asusd
}

start_pre() {
        /usr/bin/sleep 2
}
