# sys-power/asusctl/files/asus-notify
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors

description="ASUS Notifications"
command="/usr/bin/asus-notify"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
supervisor=supervise-daemon

depend() {
        after asusd
}

start_pre() {
        /usr/bin/sleep 2
}
